# Scripts/pcapworker.service
[Unit]
Description=PCAP Worker - ekstrak fitur dari pcap (PCAPWorker.py)
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Ganti 'pros' jika kamu memakai user lain; jangan pakai root kecuali perlu.
User=pros
Group=pros
WorkingDirectory=YOUR_PATH
ExecStart=YOUR_PATH/venv/bin/python YOUR_PATH/PCAPWorker.py

# Restart otomatis saat gagal
Restart=on-failure
RestartSec=5s
# Output log diarahkan ke journal (systemd journal)
StandardOutput=journal
StandardError=journal

# Opsional: batasi resources (contoh)
# LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
